# vagrant/conf/supervisord/metadata.conf
[program:metadata]
directory=/usr/local/bin
command=/usr/local/bin/metadata -c /conf/services/metadata.conf
stderr_logfile=/var/log/metadata.log
stdout_logfile=/var/log/metadata.log
autorestart=true
