# vagrant/conf/supervisord/rooms.conf
[program:rooms]
directory=/usr/local/bin
command=/usr/local/bin/rooms -c /conf/services/rooms.toml
stderr_logfile=/var/log/rooms.log
stdout_logfile=/var/log/rooms.log
autorestart=true
