# vagrant/conf/supervisord/tracking.conf
[program:tracking]
directory=/usr/local/bin
command=/usr/local/bin/tracking -c /conf/services/tracking.toml
stderr_logfile=/var/log/tracking.log
stdout_logfile=/var/log/tracking.log
autorestart=true
