# vagrant/conf/supervisord/welcome.conf
[program:welcome]
directory=/usr/local/bin
command=/usr/local/bin/welcome
stderr_logfile=/var/log/welcome.log
stdout_logfile=/var/log/welcome.log
autorestart=true
