# vagrant/conf/indexer.conf
[program:indexer]
directory=/usr/local/bin
command=/usr/local/bin/indexer
stderr_logfile=/var/log/indexer.log
stdout_logfile=/var/log/indexer.log
