# vagrant/conf/rooms.conf
[program:rooms]
directory=/usr/local/bin
command=/usr/local/bin/rooms
stderr_logfile=/var/log/rooms.log
stdout_logfile=/var/log/rooms.log
