# vagrant/conf/soapbox.conf
[program:soapbox]
directory=/usr/local/bin
command=/usr/local/bin/soapbox
stderr_logfile=/var/log/soapbox.log
stdout_logfile=/var/log/soapbox.log
user=nginx
