# vagrant/conf/voicely.conf
[program:voicely]
directory=/usr/local/bin
command=/usr/local/bin/voicely
stderr_logfile=/var/log/voicely.log
stdout_logfile=/var/log/voicely.log
